# server/deploy/text-calculator-sync.service
# systemd unit — serveur de synchronisation Text Calculator.
#
# Installation :
#   sudo cp deploy/text-calculator-sync.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now text-calculator-sync
#
# Toute la configuration vit dans le fichier « .env » (voir .env.example),
# chargé ci-dessous par EnvironmentFile — vous n'avez donc plus rien à éditer
# ici une fois les chemins ajustés à votre VPS (User, WorkingDirectory, chemin
# du .env).

[Unit]
Description=Text Calculator sync server
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/text-calculator/server
EnvironmentFile=/var/www/text-calculator/server/.env
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
